# .config/systemd/user/eye-break.timer
[Unit]
Description=Run eye break reminder every 20 minutes

[Timer]
OnCalendar=*:0/20
Persistent=true

[Install]
WantedBy=timers.target
