# .config/systemd/user/eye-break.timer
[Unit]
Description=Run break reminder every 25 minutes

[Timer]
OnBootSec=25min
OnUnitActiveSec=25min
Persistent=true

[Install]
WantedBy=timers.target
