# .config/autostart/noxrdb.desktop
[Desktop Entry]
Type=Application
Name=NoXRDB
Exec=/usr/bin/xrdb -remove
StartupNotify=false
OnlyShowIn=KDE;
